# templates/com.user.crankreview.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.crankreview</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/PATH_TO_REPO/scripts/run_weekly_review.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/PATH_TO_REPO/logs/automation.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/PATH_TO_REPO/logs/automation_error.log</string>
</dict>
</plist>